# ordo/systemd/ordo-sync.service
[Unit]
Description=Ordo Sync Daemon (watch + remote poll)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
NotifyAccess=all
WatchdogSec=60

# Tune remote pull cadence during work hours; raise to 300–600 when idle
Environment=ORDO_REMOTE_POLL_INTERVAL=120
# Keep large initial syncs unconstrained; set to a value only if needed
Environment=ORDO_BISYNC_TIMEOUT_SEC=0
# Keep excludes on for normal runs
Environment=ORDO_USE_EXCLUDES=1

ExecStart=%h/projetos/hub/OrdoMount/ordo/scripts/ordo-sync.sh daemon 300
Restart=on-failure
RestartSec=15
SuccessExitStatus=143

TimeoutStopSec=30
KillSignal=SIGTERM
KillMode=control-group

NoNewPrivileges=true
SyslogIdentifier=ordo-sync

[Install]
WantedBy=default.target
